# recipe/files/etc/systemd/system/upload-image.timer
[Unit]
Description=Timer for uploading image.

[Timer]
OnBootSec=3min
OnUnitActiveSec=60min

[Install]
WantedBy=timers.target
