# recipe/templates/etc/systemd/system/generate-time-lapse.service.erb
[Unit]
Description=Generate time-lapse movie.
RefuseManualStart=no
RefuseManualStop=yes

[Service]
Type=oneshot
User=<%= @user %>
Environment="PROCESS_SPOOL_DIRECTORY=<%= @spool_dir %>"
Environment="AWS_ACCESS_KEY_ID=<%= @aws_key_id %>"
Environment="AWS_SECRET_ACCESS_KEY=<%= @aws_secret %>"
Environment="AWS_DEFAULT_REGION=<%= @aws_region %>"
Environment="S3_BUCKET=<%= @s3_bucket %>"
Environment="VIDEO_DEVICES=<%= @video_devices %>"
WorkingDirectory=/opt/time-lapse-camera/bin
ExecStart=/opt/time-lapse-camera/bin/GENERATE_TIME_LAPSE
