# recipe/templates/etc/systemd/system/time-lapse-take-picture.service.erb
[Unit]
Description=Take picture for time-lapse movie.
RefuseManualStart=no
RefuseManualStop=yes

[Service]
Type=oneshot
User=<%= @user %>
Environment="SPOOL_DIRECTORY=<%= @spool_dir %>"
Environment="IMAGE_RESOLUTION=<%= @image_resolution %>"
WorkingDirectory=/opt/time-lapse-camera/bin
ExecStart=/opt/time-lapse-camera/bin/TAKE_PICTURE
